# etc/cron.d/update-blacklist-unbound
SHELL='/bin/bash'

PATH='/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin'

00 00 * * * root bash '/home/snwmds/AmnDNS/scripts/update-blacklist.sh'